# supervisor/reverb.conf
; Drop into /etc/supervisor/conf.d/reverb.conf (or whatever your distro uses).
; Adjust the user, paths, and PHP binary to match your server.
;
; After editing:
;   sudo supervisorctl reread && sudo supervisorctl update
;
; Graceful reload after code or config changes:
;   sudo -u forge php /home/forge/reverb/artisan reverb:restart
;
; Supervisor will auto-restart the daemon after the graceful exit.

[program:reverb]
process_name=%(program_name)s_%(process_num)02d
command=php /home/forge/reverb/artisan reverb:start --host=0.0.0.0 --port=8080
autostart=true
autorestart=true
user=forge
redirect_stderr=true
stdout_logfile=/var/log/reverb.log
stopwaitsecs=3600
numprocs=1
; Bump file descriptor limit if you expect heavy concurrency. Pair with
; /etc/security/limits.conf and supervisord's own minfds setting.
; minfds=10000
